# packaging/debian/postinst
#!/bin/sh
set -e

case "$1" in
    configure)
        if ! getent passwd helixcode >/dev/null 2>&1; then
            adduser --system --group --home /var/lib/helixcode helixcode
        fi
        mkdir -p /var/lib/helixcode /etc/helixcode
        chown -R helixcode:helixcode /var/lib/helixcode
        if [ ! -f /etc/helixcode/config.yaml ]; then
            cp /usr/share/helixcode/config.yaml.example /etc/helixcode/config.yaml
        fi
        systemctl daemon-reload
        systemctl enable helixcode.service || true
        systemctl start helixcode.service || true
        ;;
esac
